# riverlevels.cron
#!/bin/sh

RECIPIENTS='colin@pythontech.co.uk'

tmp=/tmp/riverlevels.tmp

python -m riverlevels alerts > "$tmp" 2>&1
if [ -s "$tmp" ]; then
    (
	echo "To: $RECIPIENTS"
	echo "Subject: River level changes"
	echo ""
	cat "$tmp"
    ) | /usr/sbin/sendmail -t -oi
fi
